# config/unibrowse.service
[Unit]
Description=unibrowse - Bridging AI agents and browser automation
After=network.target

[Service]
Type=simple
User=ben
WorkingDirectory=/mnt/nas/data/code/unibrowse
Environment="NODE_ENV=production"
ExecStart=/usr/bin/node --env-file=/mnt/nas/data/code/unibrowse/.env /mnt/nas/data/code/unibrowse/dist/http-server.js
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=unibrowse

# Security
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/mnt/nas/data/code/unibrowse

[Install]
WantedBy=multi-user.target
